# docker/marketplace-php-apache/crontab
* * * * * application . /get_app_root.sh; php $APP_ROOT/artisan schedule:run >> /dev/null 2>&1